# files/etc/uci-defaults/99-custom
uci -q batch <<- 'EOF'

set network.lan.ipaddr='10.0.0.1'

set luci.main.mediaurlbase='/luci-static/argon'

set wizard.default.ipv6='0'

set wireless.default_radio0.ieee80211r='1'
set wireless.default_radio1.ieee80211r='1'
EOF

root_password="root"
if [ -n "$root_password" ]; then
  (echo "$root_password"; sleep 1; echo "$root_password") | passwd > /dev/null
fi

echo "All configurations applied successfully!"