# files/etc/uci-defaults/99-custom
uci -q batch << 'EOI'
set network.lan.ipaddr='192.168.10.1'
set luci.main.mediaurlbase='/luci-static/argon'
set system.hostname='Router'
set rpcd.password='$p$root'
EOI
