# files/etc/uci-defaults/99-custom
#!/bin/bash

uci -q batch <<- 'EOF'

set network.lan.ipaddr='10.0.0.1'
set luci.main.mediaurlbase='/luci-static/argon'

set wizard.default.ipv6='0'

set wireless.radio0.mu_beamformer='1'
set wireless.radio1.mu_beamformer='1'

set wireless.radio0.country='US'
set wireless.radio1.country='US'

set wireless.radio0.htmode='HE160'

set wireless.default_radio0.ieee80211r='1'
set wireless.default_radio1.ieee80211r='1'

set wireless.radio0.channel='149'

set upnpd.config.enabled='1'

set system.@system[0].hostname='OpenWRT'

set passwall.AVVjwzzU=nodes
set passwall.AVVjwzzU.remarks='自动切换'
set passwall.AVVjwzzU.type='Socks'
set passwall.AVVjwzzU.address='127.0.0.1'
set passwall.AVVjwzzU.port='1081'

set passwall.@global[0].tcp_node='AVVjwzzU'

set passwall.oFIfqfES=socks
set passwall.oFIfqfES.enabled='1'
set passwall.oFIfqfES.port='1081'
set passwall.oFIfqfES.http_port='0'
set passwall.oFIfqfES.enable_autoswitch='1'

EOF

root_password="root"
if [ -n "$root_password" ]; then
  (echo "$root_password"; sleep 1; echo "$root_password") | passwd > /dev/null
fi

echo "All configurations applied successfully!"