# files/etc/uci-defaults/99-custom
uci -q batch <<- 'EOF'

set network.lan.ipaddr='10.0.0.1'
set luci.main.mediaurlbase='/luci-static/argon'

set wizard.default.ipv6='0'

set wireless.radio0.mu_beamformer='1'
set wireless.radio1.mu_beamformer='1'

set wireless.radio0.cell_density='3'
set wireless.radio1.cell_density='3'

set wireless.radio0.country='US'
set wireless.radio1.country='US'

uci set wireless.radio0.htmode='HE160'

EOF

root_password="root"
if [ -n "$root_password" ]; then
  (echo "$root_password"; sleep 1; echo "$root_password") | passwd > /dev/null
fi

echo "All configurations applied successfully!"