# files/etc/uci-defaults/99-custom
#!/bin/bash

uci -q batch <<- 'EOF'

set network.lan.ipaddr='10.0.0.1'
set luci.main.mediaurlbase='/luci-static/argon'

set wizard.default.ipv6='0'

set wireless.radio0.mu_beamformer='1'
set wireless.radio1.mu_beamformer='1'

set wireless.radio0.country='US'
set wireless.radio1.country='US'

set wireless.radio0.htmode='HE160'

set wireless.default_radio0.ieee80211r='1'
set wireless.default_radio1.ieee80211r='1'

set wireless.radio0.channel='149'

set upnpd.config.enabled='1'

set system.@system[0].hostname='OpenWRT'

# 添加新的节点配置并指定ID 设置该节点的各项参数
set passwall.SocksNOD=nodes
set passwall.SocksNOD.remarks='自动切换'
set passwall.SocksNOD.type='Socks'
set passwall.SocksNOD.address='127.0.0.1'
set passwall.SocksNOD.port='1081'

set passwall.global.tcp_node = 'SocksNOD'
EOF

root_password="root"
if [ -n "$root_password" ]; then
  (echo "$root_password"; sleep 1; echo "$root_password") | passwd > /dev/null
fi

echo "All configurations applied successfully!"